# ops/hermes/launchd/ai.hermes.caretaker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.hermes.caretaker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/leekaiyi/workspace/agent builder/ops/hermes/run-caretaker.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/leekaiyi/.hermes/logs/caretaker.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/leekaiyi/.hermes/logs/caretaker.error.log</string>
</dict>
</plist>
